# systemd/martiandtrump.service
[Unit]
Description=Martian Donald Trump translation bot

[Service]
Type=oneshot
ExecStart=/path/to/martiandtrump.py
User=YourUsername
